# supervisor/mqtt.conf
[program:mqtt-listener]
command=python3 /home/ubuntu/app/mqtt/mqtt_to_db.py
chown=ubuntu:supervisor